# deploy/systemd/tiny-api.service
[Unit]
Description=tiny beanstalkd enqueue API
After=network.target beanstalkd.service
Wants=beanstalkd.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/git_projects/tiny
ExecStart=/usr/local/bin/tiny-api -listen :8080 -beanstalk-addr 127.0.0.1:11300 -tube jobs -cors-origin https://tinyjobsui.elladali.com
Restart=always
RestartSec=3
Environment=TZ=UTC

# Optional hardening for production:
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/ubuntu/git_projects/tiny

[Install]
WantedBy=multi-user.target
